# bug_zapper.service
[Unit]
Description=Bug Zapper Kill Streak Tracker
After=network.target sound.target
Wants=sound.target

[Service]
ExecStart=/home/danny/.cache/pypoetry/virtualenvs/bug-zapper-EZ41eyuC-py3.12/bin/python /home/danny/bug-zapper/bug_zapper.py
WorkingDirectory=/home/danny/bug-zapper
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=5
User=danny
Group=audio
Environment=PATH=/home/danny/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=ALSA_CARD=S3
TimeoutStopSec=10
Nice=-10

[Install]
WantedBy=multi-user.target
